# deploy/systemd/runtime_service.service
[Unit]
Description=LLM Telegram Runtime Service
After=network.target

[Service]
Type=simple
User=dev
WorkingDirectory=/opt/llm/llm-telegram-client-in-dev
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/llm/llm-telegram-client-in-dev/.venv/bin/python /opt/llm/llm-telegram-client-in-dev/runtime_service.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
