# deploy/systemd/telegram_service.service
[Unit]
Description=LLM Telegram Bot Service
After=network.target api_service.service
Wants=api_service.service

[Service]
Type=simple
User=dev
WorkingDirectory=/opt/llm/llm-telegram-client-in-dev
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/llm/llm-telegram-client-in-dev/.venv/bin/python /opt/llm/llm-telegram-client-in-dev/telegram_service.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
